# phosphor-ipmi-kcs@.service.in
[Unit]
Description=Phosphor IPMI KCS DBus Bridge
Wants=phosphor-ipmi-host.service
After=phosphor-ipmi-host.service

[Service]
Type=notify
Restart=always
ExecStart=@BIN@ -c "%i" -t 0
SyslogIdentifier=kcsbridged-%i

[Install]
WantedBy=multi-user.target
RequiredBy=
